# x-ui.service
[Unit]
Description=x-ui Service
After=network.target
Wants=network.target

[Service]
Type=simple
WorkingDirectory=/usr/local/x-ui-0.0.3
ExecStart=/usr/local/x-ui-0.0.3

[Install]
WantedBy=multi-user.target